# systemd/drop-server.service
[Unit]
Description=Nest file drop server (HTTPS upload to ~/Desktop/Nest)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/github/.willow/drop.env
Environment=WILLOW_ROOT=%h/github/willow-2.0
Environment=PYTHONPATH=%h/github/willow-2.0
ExecStart=%h/github/willow-2.0/willow.sh exec-python %h/github/willow-2.0/tools/drop_server.py
StandardOutput=append:/tmp/drop-server.log
StandardError=append:/tmp/drop-server.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
