# systemd/grove-serve.service
[Unit]
Description=Grove HTTP server — REST API + Willow/FRANK watch loop
Documentation=https://github.com/rudi193-cmd/willow-2.0
After=network.target postgresql.service

[Service]
Type=simple
Environment=WILLOW_PG_DB=willow_20
Environment=WILLOW_PG_USER=%u
Environment=WILLOW_ROOT=%h/github/willow-2.0
Environment=WILLOW_GROVE_PORT=7777
Environment=PYTHONPATH=%h/github/willow-2.0
ExecStart=%h/github/willow-2.0/willow.sh exec-python -m core.grove_serve --port 7777
WorkingDirectory=%h/github/willow-2.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
