# systemd/host-divergence-watch.timer
[Unit]
Description=Willow Host-Divergence Watchdog — daily pass
Documentation=https://github.com/rudi193-cmd/willow-2.0

[Timer]
OnCalendar=daily
OnBootSec=15min
Persistent=true
Unit=host-divergence-watch.service

[Install]
WantedBy=timers.target
