# systemd/orin-worker.service
[Unit]
Description=Orin — mistral:7b batch processor sub-agent
After=network.target postgresql.service willow-metabolic.socket

[Service]
Type=simple
Environment=WILLOW_AGENT_NAME=orin
EnvironmentFile=-%h/github/.willow/env
ExecStart=%h/github/willow-2.0/willow.sh exec-python -m agents.orin.worker
WorkingDirectory=%h/github/willow-2.0
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orin-worker

[Install]
WantedBy=default.target
