# systemd/repo-fleet-sweep.service
[Unit]
Description=Willow Repo Fleet Sweep — read-only hygiene report across ~/github
Documentation=https://github.com/rudi193-cmd/willow-2.0

[Service]
Type=oneshot
Environment=WILLOW_ROOT=%h/github/willow-2.0
ExecStart=/usr/bin/python3 %h/github/willow-2.0/scripts/repo_fleet_sweep.py --root %h/github --emit-flags

[Install]
WantedBy=default.target
