# systemd/willow-metabolic.service
[Unit]
Description=Willow Metabolic Service — Norn Scheduler
Documentation=https://github.com/rudi193-cmd/willow-2.0

[Service]
Type=oneshot
Environment=WILLOW_ROOT=%h/github/willow-2.0
Environment=WILLOW_STORE_ROOT=%h/.willow/store
Environment=WILLOW_PG_DB=willow_20
ExecStart=/bin/sh -c 'python3 "$WILLOW_ROOT/core/metabolic.py"'

[Install]
WantedBy=default.target
