# systemd/willow-w8-census.service
[Unit]
Description=Willow W8 Census Witness — canonical reconstruction heartbeat
Documentation=https://github.com/rudi193-cmd/willow-2.0

[Service]
Type=oneshot
Environment=WILLOW_ROOT=%h/github/willow-2.0
Environment=WILLOW_STORE_ROOT=%h/.willow/store
Environment=WILLOW_PG_DB=willow_20
ExecStart=%h/github/willow-2.0/willow.sh w8-census

[Install]
WantedBy=default.target
