# systemd/willow-grove-listen.service
[Unit]
Description=Willow Grove LISTEN/NOTIFY mention monitor (grove_listen.py)
Documentation=https://github.com/rudi193-cmd/willow-1.9
After=postgresql.service network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=WILLOW_PG_DB=willow_20
Environment=WILLOW_PG_USER=%u
Environment=WILLOW_ROOT=%h/willow-1.9
Environment=PYTHONPATH=%h/willow-1.9
Environment=WILLOW_AGENT_NAME=Auto
# Direct: @Auto (and aliases in ALIASES). Broadcasts: @all / @All (case-insensitive).
# Optional additive identities: Environment=GROVE_MENTION_WATCH=hanuman,loki
# When unset and primary≠Auto, grove_listen.py also watches Auto by default.
ExecStart=/usr/bin/python3 %h/willow-1.9/willow/grove_listen.py
StandardOutput=append:/tmp/grove-monitor.log
StandardError=append:/tmp/grove-monitor.log
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
