# systemd/willow-mcp.service
[Unit]
Description=Willow SAP MCP Server — HTTP (SSE) transport for IDE clients
Documentation=https://github.com/rudi193-cmd/willow-1.9
After=network.target postgresql.service

[Service]
Type=simple
Environment=WILLOW_ROOT=%h/willow-1.9
Environment=WILLOW_PG_DB=willow_19
Environment=WILLOW_PG_USER=%u
Environment=PYTHONPATH=%h/willow-1.9
ExecStart=%h/willow-1.9/.venv-dev/bin/python3 %h/willow-1.9/sap/sap_mcp.py --http
WorkingDirectory=%h/willow-1.9
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
