# apps/ratatosk/systemd/ratatosk-listener.service
[Unit]
Description=Ratatosk desktop Grove listener
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/github/willow-2.0/apps/ratatosk
Environment=RATATOSK_TRANSPORT=tailnet
ExecStart=%h/github/willow-2.0/.venv/bin/ratatosk listen --channel dispatch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
